# deploy/fast-frame.service
[Unit]
Description=Fast-Frame - AI API Gateway Platform
Documentation=https://github.com/meitianwang/fast-frame
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=fast-frame
Group=fast-frame
WorkingDirectory=/opt/fast-frame
ExecStart=/opt/fast-frame/fast-frame
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fast-frame

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/fast-frame

# Environment - Server configuration
# Modify these values to change listen address and port
Environment=GIN_MODE=release
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080

[Install]
WantedBy=multi-user.target
